# etc/systemd/system/container-caddy-certificates@.path
[Unit]
Description=Caddy certificates monitor
PartOf=%N.service

[Path]
PathChanged=/mnt/data/var/lib/caddy/certificates/acme-staging-v02.api.letsencrypt.org-directory/%i/%i.json
PathChanged=/mnt/data/var/lib/caddy/certificates/acme-v02.api.letsencrypt.org-directory/%i/%i.json
PathChanged=/mnt/data/var/lib/caddy/certificates/acme.zerossl.com-v2-DV90/%i/%i.json

[Install]
WantedBy=multi-user.target
