# etc/systemd/system/ssh-key-dir.path
[Unit]
Description=SSH authorized key directory monitor
PartOf=ssh-key-dir.service

[Path]
PathChanged=/mnt/data/.ssh/authorized_keys.d
MakeDirectory=yes

[Install]
WantedBy=multi-user.target
