# etc/systemd/system/ssh-key-dir.path
[Unit]
Description=SSH authorized key directory monitor
PartOf=%N.service

[Path]
PathChanged=/mnt/data/etc/ssh/authorized_keys.d
MakeDirectory=yes

[Install]
WantedBy=multi-user.target
